# wiki.service
[Unit]
Description=local wiki server

[Service]
ExecStart=/usr/bin/wiki

[Install]
WantedBy=multi-user.target
